# telemt-panel.service
[Unit]
Description=Telemt Panel
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/telemt-panel --config /etc/telemt-panel/config.toml
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Hardening (compatible with update and config editing features)
NoNewPrivileges=true
ProtectHome=true

[Install]
WantedBy=multi-user.target
